# infra/crawler/crontab
# biblioHack autonomous crawl schedule (UTC) — interpreted by supercronic.
# Logs to stdout: `docker logs bibliohack-crawler`.
#
# A shared flock (run-job.sh) serialises the two jobs. The refresh runs at :30
# (NOT :00) specifically so it never fires at the same minute as the 03:00
# nightly — otherwise they race for the lock and one skips ("another crawl job
# is still running"), which previously starved the nightly every day.
#
# Catalogue growth every hour at :00 — discover advances a resumable cursor
# through the result set, then the worker ingests the new slice. With
# DISCOVER_MAX/WORKER_MAX ~200 this finishes well before :30.
0 * * * * /app/run-job.sh discover_worker
# Availability refresh every 6 hours at :30 (offset from the hourly growth so
# they never race for the flock). Re-scrapes records whose probe is due.
30 */6 * * * /app/run-job.sh refresh
# Cover resolution every 2 hours at :15 — off the OPAC (hits Open Library /
# Google Books), own lock, fills the shared cover store for new ISBNs.
15 */2 * * * /app/run-job.sh covers
# Embeddings every 3 hours at :45 — off the OPAC (hits the HF Inference API),
# own lock, vectorizes records lacking an embedding for semantic search.
45 */3 * * * /app/run-job.sh embed
# Catalogue relevance nightly at 04:00 — pure DB compute (no OPAC, no network),
# own lock, rescores every record from the availability series + holdings so
# /browse and search lead with the best titles. Offset from the 03:00 growth
# nightly so it reads a settled snapshot set.
0 4 * * * /app/run-job.sh relevance
