# infra/crawler/crontab
# biblioHack autonomous crawl schedule (UTC) — interpreted by supercronic.
# Logs to stdout: `docker logs bibliohack-crawler`.
#
# A shared flock (run-job.sh) serialises the two jobs. The refresh runs at :30
# (NOT :00) specifically so it never fires at the same minute as the 03:00
# nightly — otherwise they race for the lock and one skips ("another crawl job
# is still running"), which previously starved the nightly every day.
#
# Catalogue growth every hour at :00 — discover advances a resumable cursor
# through the result set, then the worker ingests the new slice. With
# DISCOVER_MAX/WORKER_MAX ~200 this finishes well before :30.
0 * * * * /app/run-job.sh discover_worker
# Availability refresh every 6 hours at :30 (offset from the hourly growth so
# they never race for the flock). Re-scrapes records whose probe is due.
30 */6 * * * /app/run-job.sh refresh
